# service/fail2ban.service.d/override.conf
[Service]
ExecStartPost=/etc/fail2ban/scripts/notify_startup.sh
